# config/systemd/polybar-i3bar.service
[Unit]
Description=Polybar

[Service]
Type=simple
ExecStart=/usr/bin/polybar -c %h/.config/polybar/config i3bar
